# deploy/elf2-server-web.service
[Unit]
Description=ELF2 training and board-dispatch web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=laipengxu
WorkingDirectory=/home/laipengxu/server
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/laipengxu/server/serve.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
